# scripts/systemd/prometeu-gateway.service
[Unit]
Description=Prometeu Gateway (FastAPI)
After=network-online.target llama-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/prometeu
Environment=PROMETEU_CONFIG=/opt/prometeu/config.json
Environment=PROMETEU_WEB_DIR=/opt/prometeu/web
Environment=PROMETEU_LLAMA_URL=http://127.0.0.1:8080
ExecStart=/opt/prometeu/venv/bin/uvicorn gateway.app:app --host 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
